# general_data_aggregator/crontab
* * * * * python /opt/app/src/main.py >> /logfile
